# ci/run-and-update-service/sig-update.service
[Unit]
Description=Runs sig-update command.

[Service]
ExecStart=/usr/local/bin/sig-update
Restart=on-failure
